# main/superd-config/services/bluez.mpris-proxy.service
[Unit]
Description=mpris-proxy daemon

[Service]
Type=simple
ExecStart=/usr/bin/mpris-proxy
Restart=always
