# main/superd-config/services/calls.service
[Unit]
Description=A phone dialer and call handler (daemon mode)

[Service]
Type=simple
ExecStart=/usr/bin/gnome-calls --daemon
Restart=always
