# main/superd-config/services/chatty.service
[Unit]
Description=Chats (daemon)

[Service]
Type=simple
ExecStart=/usr/bin/chatty --daemon
Restart=always
