# main/superd-config/services/mmsd-tng.service
[Unit]
Description=Daemon for MMS

[Service]
Type=simple
ExecStart=/usr/bin/mmsdtng -d
Restart=always
