# main/superd-config/services/vvmd.service
[Unit]
Description=Daemon for Visual Voicemail

[Service]
Type=simple
ExecStart=/usr/bin/vvmd -d
Restart=always
